# supervisord.conf
[program:php5-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user = root
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
